# docker-saltCert.service
[Unit]
Description=saltCert plugin
Before=docker.service

[Service]
ExecStart=/usr/bin/docker-saltCert

[Install]
WantedBy=docker.service
